# deploy/systemd/xsignal-vpe-testnet-auto-cycle.service
[Unit]
Description=X-Signal VPE testnet automatic strategy cycle
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/x-signal
EnvironmentFile=/etc/xsignal/binance-testnet.env
ExecStart=/opt/x-signal/.venv/bin/xsignal-vpe-live run-cycle --mode testnet --db /var/lib/xsignal/live/vpe-testnet.sqlite --lookback-bars 120
